# rootfs/etc/systemd/system/power-profiles-daemon.service.d/10-sp11-platform-profile.conf
[Service]
ExecStart=
ExecStart=/usr/local/libexec/power-profiles-daemon-sp11
