# rootfs/etc/systemd/system/sp11-noidle.service
[Unit]
Description=SP11 stability: block deep CPU idle (state1) - RPMh sleep-set fabric wedge mitigation (2026-07-17)
# Root cause under investigation: deep-idle entry applies RPMh sleep-set
# votes; racing a display transition starves scan-out and wedges the NoC
# (suspend hang, blacklist-boot lockups, spontaneous GDM freezes).
# Cost: higher idle power (WFI only). Remove once a real fix lands.
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/sh -ec 'if grep -qw cpuidle.off=1 /proc/cmdline; then exit 0; fi; count=0; for f in /sys/devices/system/cpu/cpu*/cpuidle/state1/disable; do [ -e "$f" ] || continue; echo 1 > "$f"; [ "$(cat "$f")" = 1 ]; count=$((count + 1)); done; [ "$count" -eq 12 ] || { echo "expected 12 effective state1 controls, found $count" >&2; exit 1; }'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
